# openhms-800.service
[Unit]
Description=OpenHMS-800 Unified Service
After=network.target

[Service]
User=SEARCH_USER
WorkingDirectory=SEARCH_DIR
ExecStart=SEARCH_DIR/.venv/bin/python3 -m openhms800.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
